# mumble.service
[Unit]
Description=mumble - whisper.cpp dictation daemon
After=sound.target

[Service]
Type=simple
WorkingDirectory=%h/projects/mumble
ExecStart=%h/projects/mumble/.venv/bin/python %h/projects/mumble/whisper_daemon.py --config %h/projects/mumble/config.toml
Restart=always
RestartSec=5

# Wayland environment
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
